# welovewe/conf/supervisord/nginx-program.conf
[program:nginx]
process_name = nginx
command = /opt/nginx/sbin/nginx
